# roles/result/files/result.service
[Unit]
Description = <project description>
After=network.target

[Service]
User=user5
EnvironmentFile=/etc/environment
WorkingDirectory= /home/result
ExecStart=/bin/node server.js
Restart=always

[Install]
WantedBy=multi-user.target